# config/borgmatic/borgmatic.cron
# Run nightly backup.
0 3 * * * root PATH=$PATH:/usr/bin:/usr/local/bin /usr/local/bin/borgmatic --verbosity 1 > /evix/logs/backup.log 2>&1
